# supervisor/conf.d/murmurd.conf
[program:murmurd]
command=/usr/sbin/murmurd -fg -ini /mumble/mumble-server.ini
user=mumble-server

